# model-monitor-service/systemd/model-monitor-service.service
[Unit]
Description=Model Monitor Service
Documentation=https://github.com/TheGreatGooo/HomeLabSmith
After=network.target
Requires=network.target

[Service]
Type=simple
User=home-lab-services
Group=home-lab-services
WorkingDirectory=/config/HomeLabSmith/model-monitor-service
Restart=always
RestartSec=10
Environment=PYTHONPATH=/config/HomeLabSmith
Environment=FLASK_ENV=production
# Executable path - using Gunicorn for production
ExecStart=/var/lib/home-lab-services/model-monitor-venv/bin/gunicorn --bind 0.0.0.0:5003 --workers 2 --timeout 120 --keep-alive 5 --max-requests 1000 --max-requests-jitter 100 --access-logfile - --error-logfile - model_monitor:app
# Logging configuration
StandardOutput=journal
StandardError=journal
SyslogIdentifier=model-monitor-service
# Resource limits
LimitNOFILE=65536
LimitNPROC=4096
# Security settings
PrivateTmp=true
ProtectSystem=strict
ReadWritePaths=/config/HomeLabSmith/model-monitor-service

[Install]
WantedBy=multi-user.target