# inference-service/systemd/inference-service.service
[Unit]
Description=Inference Model Service
Documentation=https://github.com/TheGreatGooo/HomeLabSmith
After=network.target
Requires=network.target

[Service]
Type=simple
User=home-lab-services
Group=home-lab-services
WorkingDirectory=/config/HomeLabSmith/inference-service
Restart=always
RestartSec=10
Environment=PYTHONPATH=/config/HomeLabSmith
Environment=FLASK_ENV=production
Environment=MODELS_CONFIG_DIR=/config/models
# Executable path - using Gunicorn for production
ExecStart=/var/lib/home-lab-services/inference-venv/bin/gunicorn --bind 0.0.0.0:5002 --workers 1 --timeout 120 --keep-alive 5 --max-requests 1000 --max-requests-jitter 100 inference_service:app
# Logging configuration
StandardOutput=journal
StandardError=journal
SyslogIdentifier=inference-service
# Resource limits
LimitNOFILE=65536
LimitNPROC=4096
# Security settings
PrivateTmp=true
ProtectSystem=strict
ReadWritePaths=/config/HomeLabSmith/inference-service

[Install]
WantedBy=multi-user.target