# shutdown-service/systemd/shutdown-service.service
[Unit]
Description=Shutdown Service
Documentation=https://github.com/your-repo/home-lab-services
After=network.target
Requires=network.target

[Service]
Type=simple
User=abc
Group=abc
WorkingDirectory=/config/HomeLabSmith/shutdown-service
Restart=always
RestartSec=10
Environment=PYTHONPATH=/config/HomeLabSmith
Environment=FLASK_ENV=production
# Executable path
ExecStart=/usr/bin/python3 /config/HomeLabSmith/shutdown-service/shutdown_service.py
# Logging configuration
StandardOutput=journal
StandardError=journal
SyslogIdentifier=shutdown-service
# Resource limits
LimitNOFILE=65536
LimitNPROC=4096
# Security settings
NoNewPrivileges=true
PrivateTmp=true
ProtectSystem=strict
ReadWritePaths=/config/HomeLabSmith/shutdown-service
# User permissions for shutdown commands
User=root
Group=root

[Install]
WantedBy=multi-user.target